# crontab
* * * * * docker ps -f health=unhealthy --format "docker restart {{.ID}}" | sh
